# arch-configs/systemd/user/battery.service
[Unit]
Description=Show warning on low battery

[Service]
;Type=oneshot
ExecStart=/home/cousine/bin/lowbattery.sh
